# examples/service/resources/debian/postrm
#!/bin/sh
# postrm script for gogenexample

set -e

update-rc.d -f gogenexample remove

#DEBHELPER#

exit 0
