# examples/service/resources/debian/postinst
#!/bin/sh
# postinst script for gogenexample

set -e

update-rc.d -f gogenexample defaults

#DEBHELPER#

exit 0
