# src/app/resources/debian/postinst
#!/bin/sh
# postinst script for ~#PROJECT#~

set -e

update-rc.d -f ~#PROJECT#~ defaults

#DEBHELPER#

exit 0
